# apps/password-pdf-generator/deploy/password-pdf-generator.service.example
[Unit]
Description=Password PDF Generator API
After=network.target

[Service]
User=passwordpdf
Group=passwordpdf
WorkingDirectory=/opt/password-pdf-generator
EnvironmentFile=/etc/password-pdf-generator.env
Environment=WIFI_PDF_CONFIG_PATH=/etc/password-pdf-generator/brand_settings.json
Environment=PATH=/opt/password-pdf-generator/.venv/bin
ExecStart=/opt/password-pdf-generator/.venv/bin/uvicorn wifi_pdf.api:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

