# apps/workflow-api/deploy/site-workflow-01.service.example
[Unit]
Description=Site Workflow 01
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=siteworkflow
Group=siteworkflow
WorkingDirectory=/opt/site-workflow-01
EnvironmentFile=/etc/site-workflow-01.env
Environment=PATH=/opt/site-workflow-01/.venv/bin
ExecStart=/opt/site-workflow-01/.venv/bin/uvicorn workflow.api:app --host 127.0.0.1 --port 8100
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
