# apps/omada-site-service/deploy/linux/omada-site-creator.service
[Unit]
Description=Omada Site Creator
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=omada-site-creator
Group=omada-site-creator
WorkingDirectory=/opt/omada-site-creator
EnvironmentFile=/etc/omada-site-creator.env
Environment=HOME=/var/lib/omada-site-creator
ExecStart=/usr/bin/npm run serve:dist
Restart=always
RestartSec=5
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
